# debian/ngcp-sems.preinst
#!/bin/sh
set -e

if [ -x "/etc/init.d/sems" ]; then
	update-rc.d sems remove >/dev/null
fi

dpkg-maintscript-helper rm_conffile /etc/default/sems -- "$@"
dpkg-maintscript-helper rm_conffile /etc/init.d/sems -- "$@"

#DEBHELPER#

exit 0
